# system_backup.timer
[Unit]
Description=System backup with Borg Timer

[Timer]
OnCalendar=Daily
Persistent=true

[Install]
WantedBy=timers.target