# system_backup.service
[Unit]
Description=Backup system with Borg

[Service]
Type=simple
ExecStart=/root/system_backup.sh

[Install]
RequiredBy=multi-user.target